# contrib/init/giantd.init
#!/bin/bash
#
#  coffeecoind The COFFEECOIN core server.
#
#
# chkconfig: 345 80 20
# description: coffeecoind
# processname: coffeecoind
#

# Source function library.
. /etc/init.d/functions

# you can override defaults in /etc/sysconfig/coffeecoind, see below
if [ -f /etc/sysconfig/coffeecoind ]; then
        . /etc/sysconfig/coffeecoind
fi

RETVAL=0

prog=coffeecoind
# you can override the lockfile via BITCOIND_LOCKFILE in /etc/sysconfig/coffeecoind
lockfile=${BITCOIND_LOCKFILE-/var/lock/subsys/coffeecoind}

# coffeecoind defaults to /usr/bin/coffeecoind, override with BITCOIND_BIN
bitcoind=${BITCOIND_BIN-/usr/bin/coffeecoind}

# coffeecoind opts default to -disablewallet, override with BITCOIND_OPTS
bitcoind_opts=${BITCOIND_OPTS}

start() {
    echo -n $"Starting $prog: "
    daemon $DAEMONOPTS $bitcoind $bitcoind_opts
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch $lockfile
    return $RETVAL
}

stop() {
    echo -n $"Stopping $prog: "
    killproc $prog
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f $lockfile
    return $RETVAL
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status $prog
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: service $prog {start|stop|status|restart}"
        exit 1
        ;;
esac
